# rolesbot.service
[Unit]
Description=RolesBot v3 – Multi-Guild Discord Bot
After=network.target

[Service]
Type=simple
User=rolesbot
WorkingDirectory=/opt/rolesbot
ExecStart=/usr/bin/java -jar /opt/rolesbot/rolesbot-3.0.0.jar
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Token via Environment Variable – NICHT im Code!
EnvironmentFile=/opt/rolesbot/.env

[Install]
WantedBy=multi-user.target
